# bin/crontab
MAILTO=evyncke@cisco.com
SHELL=/bin/bash
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
#
#  SDO Analysis to fetch from IETF, YangModels, ...
#
5 0 * * 0 (cd ~/sdo_analysis/bin/ ; ./cronjob-weekly) > /var/yang/logs/cronjob-weekly.log
5 0 * * 1-6 (cd ~/sdo_analysis/bin/ ; ./cronjob-daily) > /var/yang/logs/cronjob-daily.log
25 6 * * * (cd ~/sdo_analysis/bin/ ;  DJANGO_SETTINGS_MODULE=yangsuite.settings.develop MEDIA_ROOT=/var/yang/yangsuite ./create_yangsuite_repo.py) > /var/yang/logs/create_yangsuite_repo.log 2>&1
